# hp-ctl.service
[Unit]
Description=HP Control - Panasonic Aquarea Heat Pump Controller
Documentation=https://github.com/jonas-rem/hp-ctl
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.config/hp-ctl
ExecStart=%h/.local/bin/hp-ctl
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=default.target
